# systemd/olx-monitor.service
[Unit]
Description=OLX Monitor - alertas em tempo real de anuncios no Telegram
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=olxmonitor
Group=olxmonitor
WorkingDirectory=/opt/olx-monitor
EnvironmentFile=/opt/olx-monitor/.env
ExecStart=/opt/olx-monitor/.venv/bin/python /opt/olx-monitor/run.py \
    --config /opt/olx-monitor/monitores.yaml \
    --db /opt/olx-monitor/olx_monitor.db
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Hardening básico (opcional, remova se causar problema no seu SO)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
